# mjpg-streamer.files
plugins/input.h
plugins/output.h
mjpg_streamer.c
mjpg_streamer.h
utils.c
utils.h
ipkg-ramips_24kec/mjpg-streamer/CONTROL/conffiles
ipkg-ramips_24kec/mjpg-streamer/CONTROL/control
ipkg-ramips_24kec/mjpg-streamer/CONTROL/postinst
ipkg-ramips_24kec/mjpg-streamer/CONTROL/prerm
ipkg-ramips_24kec/mjpg-streamer/etc/config/mjpg-streamer
ipkg-ramips_24kec/mjpg-streamer/etc/hotplug.d/usb/20-mjpg-streamer
ipkg-ramips_24kec/mjpg-streamer/etc/init.d/mjpg-streamer
ipkg-ramips_24kec/mjpg-streamer/usr/bin/mjpg_streamer
ipkg-ramips_24kec/mjpg-streamer/usr/lib/input_uvc.so
ipkg-ramips_24kec/mjpg-streamer/usr/lib/output_http.so
ipkg-ramips_24kec/mjpg-streamer/www/webcam/index.html
plugins/input_control/dynctrl.c
plugins/input_control/dynctrl.h
plugins/input_control/input_uvc.c
plugins/input_control/uvc_compat.h
plugins/input_control/uvcvideo.h
plugins/input_file/input_file.c
plugins/input_file/input_file.so
plugins/input_http/input_http.c
plugins/input_http/misc.c
plugins/input_http/misc.h
plugins/input_http/mjpg-proxy.c
plugins/input_http/mjpg-proxy.h
plugins/input_http/version.h
plugins/input_ptp2/input_ptp2.c
plugins/input_ptp2/input_ptp2.h
plugins/input_testpicture/pictures/160x120_1.jpg
plugins/input_testpicture/pictures/160x120_2.jpg
plugins/input_testpicture/pictures/320x240_1.jpg
plugins/input_testpicture/pictures/320x240_2.jpg
plugins/input_testpicture/pictures/640x480_1.jpg
plugins/input_testpicture/pictures/640x480_2.jpg
plugins/input_testpicture/pictures/960x720_1.jpg
plugins/input_testpicture/pictures/960x720_2.jpg
plugins/input_testpicture/input_testpicture.c
plugins/input_testpicture/testpictures.h
plugins/input_uvc/dynctrl-logitech.h
plugins/input_uvc/dynctrl.c
plugins/input_uvc/dynctrl.h
plugins/input_uvc/huffman.h
plugins/input_uvc/input_uvc.c
plugins/input_uvc/jpeg_utils.c
plugins/input_uvc/jpeg_utils.h
plugins/input_uvc/uvcvideo.h
plugins/input_uvc/v4l2uvc.c
plugins/input_uvc/v4l2uvc.h
plugins/output_autofocus/output_autofocus.c
plugins/output_autofocus/processJPEG_onlyCenter.c
plugins/output_autofocus/processJPEG_onlyCenter.h
plugins/output_file/examples/change_filename.sh
plugins/output_file/examples/ftp_upload.sh
plugins/output_file/examples/show_filename.sh
plugins/output_file/output_file.c
plugins/output_file/output_file.h
plugins/output_http/httpd.c
plugins/output_http/httpd.h
plugins/output_http/output_http.c
plugins/output_rtsp/output_rtsp.c
plugins/output_rtsp/rtsp.c
plugins/output_tcp/examples/change_filename.sh
plugins/output_tcp/examples/ftp_upload.sh
plugins/output_tcp/examples/show_filename.sh
plugins/output_tcp/output_tcp.c
plugins/output_tcp/output_tcp.h
plugins/output_udp/output_udp.c
plugins/output_viewer/output_viewer.c
plugins/input.h
plugins/output.h
scripts/make_deb.sh
www/bodybg.gif
www/cambozola.jar
www/control.html
www/control_old.html
www/example.jpg
www/favicon.ico
www/favicon.png
www/fix.css
www/functions.js
www/go-down.png
www/go-next.png
www/go-previous.png
www/go-up.png
www/index.html
www/java.html
www/java_control.html
www/java_simple.html
www/javascript.html
www/javascript_motiondetection.html
www/javascript_simple.html
www/jquery.js
www/jquery.rotate.js
www/jquery.ui.core.min.js
www/jquery.ui.custom.css
www/jquery.ui.tabs.min.js
www/jquery.ui.widget.min.js
www/JQuerySpinBtn.css
www/JQuerySpinBtn.js
www/LICENSE.txt
www/rotateicons.png
www/sidebarbg.gif
www/spinbtn_updn.gif
www/static.html
www/static_simple.html
www/stream.html
www/stream_simple.html
www/style.css
www/videolan.html
CHANGELOG
LICENSE
mjpg_streamer.c
mjpg_streamer.h
README
README-UNSTABLE
start.sh
TODO
utils.c
utils.h
plugins/output_tcp/clienttcp.h
plugins/output_tcp/clienttcp.c
